# udev/etc/udev/rules.d/69-hdparm.rules
#ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/usr/bin/hdparm -B 254 -S 0 /dev/sdb"
