# udev/etc/udev/rules.d/99-keyboard-setup.rules
SUBSYSTEM=="input", ACTION=="add", ENV{ID_INPUT_KEYBOARD}=="1", RUN+="/home/alx/.scripts/keyboard_setup_udev.sh"
